# supervisord.conf
[supervisord]
nodaemon=true

[program:openvpn]
command=openvpn /etc/openvpn/udp1194.conf >> /etc/openvpn/udp1194.log
autostart=true
autorestart=true