# supervisord.conf
[supervisord]
nodaemon=true

[program:openvpn]
command=openvpn /etc/openvpn/udp1194.conf
autostart=true
autorestart=true